# setup/rc.local
#!/bin/sh -e
#
# rc.local
#

for i in /sys/bus/usb/devices/*/power/autosuspend; do echo -1 > $i; done

exit 0